# systemd/quick-vless-check.service
[Unit]
Description=Quick-VLESS Check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/quick-vless check
