# systemd/quick-node-check.timer
[Unit]
Description=Quick-Node periodic check

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
