# systemd/quick-vless-serve.service
[Unit]
Description=Quick-VLESS Subscription Server
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/quick-vless serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
