# systemd/quick-node-check.service
[Unit]
Description=Quick-Node Check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/quick-node check
